# deploy/mw-aging.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.marrow.aging</string>

  <key>ProgramArguments</key>
  <array>
    <string>__VENV_PYTHON__</string>
    <string>-m</string>
    <string>marrow.aging</string>
    <string>--apply</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__PROJECT_DIR__</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>__PATH_ENV__</string>
  </dict>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key><integer>0</integer>
    <key>Hour</key><integer>12</integer>
    <key>Minute</key><integer>0</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>ProcessType</key>
  <string>Background</string>

  <key>StandardOutPath</key>
  <string>__LOG_DIR__/mw-aging.log</string>
  <key>StandardErrorPath</key>
  <string>__LOG_DIR__/mw-aging.log</string>
</dict>
</plist>
