# deploy/mw-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.marrow.watcher</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/Gabrielle/cc-lab/marrow/.venv/bin/python</string>
    <string>-m</string>
    <string>marrow.watcher</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/Gabrielle/cc-lab/marrow</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/Gabrielle/.local/bin:/Users/Gabrielle/cc-lab/marrow/.venv/bin:/usr/bin:/bin</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>ThrottleInterval</key>
  <integer>5</integer>

  <key>ProcessType</key>
  <string>Background</string>

  <key>StandardOutPath</key>
  <string>/Users/Gabrielle/.config/marrow/logs/watcher-launchd.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/Gabrielle/.config/marrow/logs/watcher-launchd.log</string>
</dict>
</plist>
